# signer-vm/ops/dmj-pdfsigner.service
[Unit]
Description=DMJ PDF Signer (PDFBox + Javalin)
After=network-online.target nginx.service
Wants=nginx.service

[Service]
Type=simple
User=dmj
Group=dmj
EnvironmentFile=/etc/default/dmj-pdfsigner
WorkingDirectory=/opt/dmj/pdf-signer
ExecStart=/usr/bin/java -jar /opt/dmj/pdf-signer/pdf-signer-1.0.0.jar
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
